# roles/service_adguard/files/update_adguard.service
[Unit]
Description=test job

[Service]
Type=oneshot
ExecStart=/bin/bash /home/olivia/update_adguard.sh