# roles/service_adguard/files/update_adguard.timer
[Unit]
Description=adguard home update

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target